# deploy/linux/smart-center-local-model-training.timer
[Unit]
Description=Run Smart Center Local Model Training Export Daily

[Timer]
OnCalendar=*-*-* 04:05:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
